# contrib/unifi.service
[Unit]
Description=Unifi Controller
After=docker.service
BindsTo=docker.service

[Service]
TimeoutStartSec=0
TimeoutStopSec=30
Restart=always
RestartSec=10

WorkingDirectory=/opt/unifi

ExecStartPre=-/usr/local/bin/docker-compose kill
ExecStartPre=-/usr/local/bin/docker-compose rm -f

ExecStart=/usr/local/bin/docker-compose up --force-recreate

ExecStop=/usr/local/bin/docker-compose stop

[Install]
WantedBy=multi-user.target
